# misc/docker-cron.service
[Unit]
Description=Docker Cron

[Service]
ExecStart=/usr/local/bin/docker-cron

[Install]
WantedBy=default.target
